# packaging/alertu-daemon.service
[Unit]
Description=AlertU privileged daemon (evdev remote guard)
Documentation=https://github.com/systm-d/alertu
After=systemd-logind.service
Wants=systemd-logind.service

[Service]
Type=simple
ExecStart=/usr/local/bin/alertu-daemon --config /etc/alertu/config.toml
Restart=on-failure
RestartSec=2

# The daemon needs to read input devices and capture from the webcam. Running
# as a dedicated user in the input+video groups is preferred over root; it also
# needs to talk to logind over the system bus and place its socket in /run.
User=alertu
SupplementaryGroups=input video
RuntimeDirectory=alertu
RuntimeDirectoryMode=0750
StateDirectory=alertu
# Deliberately 0755, not 0750: the snapshot boundary lives one level down, on
# /var/lib/alertu/snapshots (0750, group-owned), and this directory must stay
# traversable so a --socket-group member can reach it.
StateDirectoryMode=0755

# Hardening (relax if it interferes with your setup).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alertu
PrivateTmp=true

[Install]
WantedBy=multi-user.target
